# riksnmp-timer.service.in
[Unit]
Description=Logs system statistics to the systemd journal
Wants=riksnmp.timer

[Service]
Type=oneshot
ExecStart=/usr/sbin/mailnew.sh /etc/riksnmp

[Install]
WantedBy=multi-user.target
