# riksnmp.service.in
[Unit]
Description=Mini SNMP service

[Service]
Restart=on-failure
ExecStart=/usr/sbin/riksnmp -ns $DAEMON_OPTS
StartLimitInterval=0
Type=dbus
BusName=xyz.openbmc_project.riksnmp

[Install]
WantedBy=basic.target
